# deploy/druns.upstart
# Upstart job configuration file

description   "Druns - Web application that schedule clients in a weekday agenda"
author        "Rafael Dantas Justo"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 10

script
  exec /usr/druns/bin/druns \
    /usr/druns/etc/webserver.conf
end script